# startup_automation/ptp-automation.service
# Archivo: /etc/systemd/system/ptp-automation.service
# Servicio systemd para automatización PTP

[Unit]
Description=PTP Automation Service
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ptp_init.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variables de entorno si son necesarias
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Límites de recursos
LimitNOFILE=65536

# Configuración de seguridad
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target